# geoconstats_supervisor.conf
[program:geoconstats]
command = <MY_APP_PATH>/gunicorn_start.sh
autostart=true
autorestart=true
stdout_logfile = <MY_APP_PATH>/log/geoconstats_errors.log
redirect_stderr = true
